# runner.service
[Unit]
Description=GitHub Runner
After=docker.service
#StartLimitIntervalSec=0

[Service]
Type=oneshot
User=runner
WorkingDirectory=/actions-runner
ExecStart=/usr/local/bin/runner.sh
EnvironmentFile=-/etc/runner/runner.conf

[Install]
WantedBy=multi-user.target